# app/supervisor/platinum_consumer.conf
[program:platinum_consumer]
command=python worker.py BaseConsumer AsyncConsumer
process_name=%(program_name)s
autostart=true
autorestart=true
stdout_logfile=/var/app/worker.log