# app/supervisor/platinum_consumer.conf
[program:platinum_consumer]
command=python worker.py Consumer PlatinumConsumer
process_name=%(program_name)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log 
