# app/supervisor/silver_consumer.conf
[program:silver_consumer]
command=python worker.py Consumer SilverConsumer
process_name=%(program_name)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log 
